# debian/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ] || [ "$1" = "deconfigure" ]; then
    systemctl stop rustyfoot-hmi.service || true
    systemctl disable rustyfoot-hmi.service || true
fi
